# systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --skip-login --noclear --nonewline --nohints --noissue --autologin simonkov %I $TERM
